# roles/local-nginx/templates/nginx.ini
[program:nginx]
command={{ brew_bin_path }}/nginx
autostart=yes
autorestart=yes
stopasgroup=yes
killasgroup=yes
