# roles/win311/templates/win311.ini
[program:win311]
directory={{ directory }}
command={{ user_bin }}/win311
autostart=yes
autorestart=yes
stopasgroup=yes
killasgroup=yes
